# src/snapman.service
[Unit]
Description=Snapman backup system

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /usr/bin/snapman

[Install]
WantedBy=multi-user.target
